# relogin-cron-task.cron
00 03 * * * bash /home/public/Git/neu-ip-gateway-manager/relogin.sh
